# container-files/etc/supervisor.d/redis.conf
[program:redis]
command = /usr/bin/redis-server /etc/redis.conf --port 6379
user = root
autorestart = false
stderr_logfile = NONE
stdout_logfile = NONE
